# deploy/relayer/accountabilibuddy-relayer.service
[Unit]
Description=AccountabiliBuddy Solana relayer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=accountabilibuddy
Group=accountabilibuddy
WorkingDirectory=/opt/accountabilibuddy
EnvironmentFile=/etc/accountabilibuddy/relayer.env
ExecStart=/usr/bin/npm run relayer
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/accountabilibuddy

[Install]
WantedBy=multi-user.target
